# systemd-services/ethumb.service.in
[Unit]
Description=Enlightenment Thumbnailer

[Service]
Type=dbus
BusName=org.enlightenment.Ethumb
ExecStart=@prefix@/bin/ethumbd
